# config/com.hermes.wq-dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hermes.wq-dashboard</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/anaconda3/bin/python3</string>
        <string>/Users/asdd/projects/wq-alpha-pipeline/dashboard/server.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/asdd/projects/wq-alpha-pipeline/dashboard</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <false/>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>StandardOutPath</key>
    <string>/Users/asdd/.wq_dashboard_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/asdd/.wq_dashboard_stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/anaconda3/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>PORT</key>
        <string>8767</string>
    </dict>
</dict>
</plist>
