# config/com.hermes.wq-workflow.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hermes.wq-workflow</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/anaconda3/bin/python3</string>
        <string>/Users/asdd/.hermes/scripts/wq_workflow_v2.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/asdd/.hermes/scripts</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>CrashedLast</key>
        <true/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>/Users/asdd/.wq_workflow_v2_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/asdd/.wq_workflow_v2_stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/anaconda3/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>UserName</key>
    <string>asdd</string>
</dict>
</plist>
